# myappProd.service
[Unit]
Description=Go Application
After=network.target

[Service]
ExecStart=/home/ubuntu/app
Type=simple

[Install]
WantedBy=default.target