# myapp.service
[Unit]
Description=Go Application
After=network.target

[Service]
ExecStart=/home/vagrant/app
Type=simple

[Install]
WantedBy=default.target